# tests/collectd/supervisord-collectd.conf
[supervisord]
nodaemon = true

[program:collectd]
directory=/
command=/opt/collectd/sbin/collectd -C /opt/collectd/etc/collectd.conf -f
stdout_logfile=/var/log/supervisor-%(program_name)s.log
stderr_logfile=/var/log/supervisor-%(program_name)s_error.log
autorestart = true
